# sysvinit/inittab

id:3:initdefault:

#si::sysinit:/etc/zinit start
l0:0:wait:/etc/zinit stop 
l3:3:wait:/etc/zinit start
l6:6:wait:/etc/zinit reboot

ca:12345:ctrlaltdel:/sbin/shutdown -t1 -a -r now
su:S016:once:/sbin/sulogin

1:2345:respawn:/sbin/agetty --noclear tty1 9600
2:2345:respawn:/sbin/agetty tty2 9600
